# .bash/profile.d/site
[[ -f ~/.bash_site ]] && source ~/.bash_site
