# .bash/profile.d/work
rs () {
	if [[ -z $1 ]]; then
		ssh rs
	else
	  ssh dev -t -A ssh "$@"
	fi
}
